# home/.local/bin/backup/photo-backup.service
[Unit]
Description=Бэкап фоток в облаке

[Service]
User=blacklacost
ExecStart=/bin/bash /opt/photo_backup.sh

[Install]
WantedBy=multi-user.target
